# deploy/systemd/stockmaster-daily-close.timer
[Unit]
Description=StockMaster Daily Close Timer

[Timer]
Unit=stockmaster-scheduler@daily-close.service
OnCalendar=Mon..Fri *-*-* 17:00:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
