# deploy/systemd/stockmaster-daily-audit-lite.timer
[Unit]
Description=StockMaster Daily Audit Lite Timer

[Timer]
Unit=stockmaster-scheduler@daily-audit-lite.service
OnCalendar=Tue..Sat *-*-* 05:30:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
