# deploy/systemd/stockmaster-weekly-training.timer
[Unit]
Description=StockMaster Weekly Training Candidate Timer

[Timer]
Unit=stockmaster-scheduler@weekly-training.service
OnCalendar=Sat *-*-* 03:30:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
