# deploy/systemd/stockmaster-weekly-calibration.timer
[Unit]
Description=StockMaster Weekly Calibration Timer

[Timer]
Unit=stockmaster-scheduler@weekly-calibration.service
OnCalendar=Sat *-*-* 06:30:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
