# deploy/systemd/stockmaster-evaluation.timer
[Unit]
Description=StockMaster Evaluation Timer

[Timer]
Unit=stockmaster-scheduler@evaluation.service
OnCalendar=Mon..Fri *-*-* 16:20:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
