# deploy/systemd/stockmaster-scheduler@.service
[Unit]
Description=StockMaster Scheduler Job (%i)
Requires=docker.service stockmaster-compose.service
After=docker.service stockmaster-compose.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/stockmaster/app
ExecStart=/bin/bash /opt/stockmaster/app/scripts/server/run_scheduler_job_host.sh %i
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target
