# deploy/systemd/stockmaster-news-after-close.timer
[Unit]
Description=StockMaster Post-close News Sync Timer

[Timer]
Unit=stockmaster-scheduler@news-after-close.service
OnCalendar=Mon..Fri *-*-* 16:10:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
