# deploy/systemd/stockmaster-daily-overlay-refresh.timer
[Unit]
Description=StockMaster Daily Overlay Refresh Timer

[Timer]
Unit=stockmaster-scheduler@daily-overlay-refresh.service
OnCalendar=Mon..Fri *-*-* 21:30:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
